# logduct/systemd/logductd.service
[Unit]
Description=Logduct Daemon

[Service]
StandardInput=socket
StandardOutput=journal
ExecStart=/usr/bin/logductd
Restart=always
User=logger
